# debian/routinator-minimal.routinator.service
[Unit]
Description=Routinator 3000
Documentation=man:routinator(1)
After=network.target

[Service]
ExecStart=/usr/bin/routinator --config=/etc/routinator/routinator.conf --syslog server

[Install]
WantedBy=multi-user.target


